# arm-openwrt-linux-gnueabi/target/usr/lib/opkg/info/kmod-c2k_pmu.list
/etc/init.d/pmu
/lib/firmware/pmu.elf
/lib/modules/3.2.26/c2k_pmu.ko
