# arm-openwrt-linux-gnueabi/target/usr/lib/opkg/info/vsftpd.list
/etc/vsftpd.conf
/usr/sbin/vsftpd
/etc/init.d/vsftpd
