# arm-openwrt-linux-gnueabi/target/usr/lib/opkg/info/webif.list
/www/cgi-bin/webif/status-connectivity.sh
/www/images/x.gif
/www/images/file.gif
/usr/lib/webif/webif.sh
/www/images/favnav_statusover_ico.jpg
/www/images/pkg.jpg
/www/cgi-bin/webif/dash.sh
/www/cgi-bin/webif/log-browse.sh
/usr/lib/webif/apn.csv
/www/cgi-bin/webif/reboot.sh
/www/cgi-bin/webif/graph_if_svg.sh
/www/images/reboot_bgbtn.gif
/www/images/disable.jpg
/etc/init.d/dmesgbackup
/www/images/calender.png
/usr/lib/webif/subcategories.awk
/www/cgi-bin/webif/log-read.sh
/usr/lib/webif/countrycodes.lst
/www/cgi-bin/webif/system-confman.sh
/www/images/lan_dash_94px.jpg
/www/images/error1.gif
/www/images/app.12.jpg
/www/svggraph/graph_cpu.svg
/www/images/transmit.png
/www/images/favnav_wanover_ico.jpg
/www/images/bg_grad_tile.jpg
/www/images/action_edit_no.gif
/usr/bin/webifmetabin
/www/js/balloontip.js
/www/svggraph/graph_if.svg
/usr/lib/webif/webif-mklanglist.sh
/etc/init.d/vlan
/www/images/buttonsmall_bg.gif
/www/images/service_enable.png
/www/cgi-bin/webif/system-settings.sh
/www/cgi-bin/webif/network-vlan.sh
/www/themes/clubman/mainnav_bg_inverse.png
/www/cgi-bin/webif/network-bridge.sh
/www/images/review_changes.gif
/www/js/waitbox.js
/www/js/forms.js
/www/js/pageload.js
/www/images/favnav_teleover_ico.jpg
/www/images/favnav_wan_ico.jpg
/www/images/qosflowc_dash_94px.jpg
/usr/lib/webif/form.awk
/www/images/edit.gif
/usr/lib/webif/validate.awk
/www/cgi-bin/webif/system-access-control.sh
/etc/opkg/xwrt.conf
/www/images/favnav_lanover_ico.jpg
/www/themes/clubman/webif.css
/www/images/vlan_dash_94px.jpg
/www/images/slidebutt_left.gif
/www/cgi-bin/webif/download.sh
/etc/config/vlan
/etc/init.d/crontab
/www/cgi-bin/webif/status-leases.sh
/www/images/telephony_dash_94px.jpg
/www/cgi-bin/webif/config.sh
/www/images/logo.jpg
/etc/config/cron
/www/cgi-bin/webif/network-wlan.sh
/www/images/qos_dash_94px.jpg
/www/images/upd.jpg
/www/themes/clubman/mainnav_bg.png
/www/js/scrollbox.js
/www/cgi-bin/webif/system-password.sh
/www/images/favnav_status_ico.jpg
/www/images/overlay.png
/www/images/help_icon.gif
/www/images/qosrate_dash_94px.jpg
/www/images/dir.gif
/www/images/abt.jpg
/usr/lib/webif/comgt.webif
/usr/lib/webif/switchinfo.sh
/www/cgi-bin/webif/status-interfaces.sh
/www/images/applyclear_bg_new.gif
/www/cgi-bin/webif/logout.sh
/etc/config/webif
/www/themes/clubman/logo.png
/www/images/action_edit.gif
/www/images/service_stop.png
/etc/httpd.conf
/usr/lib/webif/notes.awk
/www/images/favnav_wifiover_ico.jpg
/www/images/favnav_wifi_ico.jpg
/www/images/action_x.gif
/www/cgi-bin/webif/graphs-vnstat.sh
/www/.version
/www/images/wifi_dash_94px.jpg
/www/cgi-bin/webif/network-dhcpiface.sh
/www/cgi-bin/webif/data.sh
/www/cgi-bin/webif/graphs-cpu.sh
/www/cgi-bin/webif/status-usb.sh
/www/themes/clubman/mainnav_hover_bg_inverse.png
/www/cgi-bin/webif/log-dmesg.sh
/www/images/enable.gif
/www/images/header_bgtile.jpg
/www/images/app.2.jpg
/www/images/lock.png
/www/images/dot2.gif
/www/images/favnav_bg2.gif
/www/cgi-bin/webif/getdata.sh
/usr/lib/webif/dhcp_options.dat
/usr/lib/webif/editor.awk
/www/themes/clubman/mainnav_active_bg_inverse.gif
/www/images/down.gif
/www/cgi-bin/webif/network-wan.sh
/www/images/warning.png
/etc/init.d/webiffirewalllog
/www/cgi-bin/webif/status-diag.sh
/www/images/bridge_dash_94px.jpg
/www/webif.js
/www/images/action_ok.gif
/www/cgi-bin/webif/network-dhcp.sh
/www/images/wan_dash_94px.jpg
/www/images/action_sv.gif
/www/js/window.js
/www/images/dot.gif
/www/cgi-bin/webif/.categories
/usr/lib/webif/functions.sh
/usr/lib/webif/apply.sh
/usr/lib/webif/categories.awk
/www/colorize.js
/www/images/blkbox.jpg
/www/images/navbar_tile.jpg
/www/cgi-bin/webif/status-basic.sh
/www/images/favnav_bg3.gif
/www/cgi-bin/webif/graphs-if.sh
/usr/lib/webif/timezones.csv
/www/images/applyclear_bg.gif
/www/images/bkup.jpg
/www/cgi-bin/webif/network-nat.sh
/www/images/nat_dash_94px.jpg
/www/js/progress.js
/www/images/button_bg.gif
/www/favicon.ico
/www/cgi-bin/webif/graphs-subcategories.sh
/www/cgi-bin/webif/network-wlan-advance.sh
/usr/lib/webif/acl.awk
/www/cgi-bin/webif/network-services.sh
/www/cgi-bin/webif/log-setup.sh
/usr/lib/webif/sanitize.awk
/www/cgi-bin/webif/network-hosts.sh
/etc/init.d/timezone
/etc/init.d/webif
/etc/languages.root
/www/cgi-bin/webif/system-upgrade.sh
/www/images/favnav_lan_ico.jpg
/www/images/pixel.gif
/www/images/up.gif
/usr/lib/webif/common.awk
/www/images/apply_ico.jpg
/www/images/service_disable.png
/www/images/close.gif
/www/js/imgdepth.js
/www/images/favnav_dash_ico.jpg
/www/cgi-bin/webif/graphs-bandwidth.sh
/www/images/slidebutt_right.gif
/www/images/service_start.png
/www/images/action_x_no.gif
/www/themes/clubman/waitbox.css
/www/images/warning.gif
/www/cgi-bin/webif/network-lan.sh
/usr/lib/webif/webif-mkthemelist.sh
/www/images/clear_ico.jpg
/www/images/favnav_tele_ico.jpg
/www/images/loading.gif
/www/images/service_restart.png
/www/images/service_enabled.png
/www/index.html
/www/js/styleswitcher.js
/www/images/configure_new.gif
/www/images/enable.jpg
/www/images/service_disabled.png
/www/themes/clubman/name
/www/images/favnav_dashover_ico.jpg
/www/themes/clubman/mainnav_hover_bg.png
/www/cgi-bin/webif/network-route.sh
/etc/init.d/custom-user-startup
/www/cgi-bin/webif/graph_cpu_svg.sh
/usr/lib/webif/browser.awk
/www/themes/clubman/mainnav_active_bg.gif
