# arm-openwrt-linux-gnueabi/target/usr/lib/opkg/info/webif.conffiles
/etc/config/webif
/etc/config/cron
/etc/httpd.conf
/etc/init.d/custom-user-startup
