# staging_dir/etc/hotplug.d/iface/01-qos
# Skip fake devices (e.g. relayd)
grep -qs "^ *$DEVICE:" /proc/net/dev || exit 0

case "$ACTION" in
	ifup)
		local prio

		config_load qos

		if [ -e /sys/class/net/$DEVICE/default_priority ]; then

			config_get prio "$DEVICE" priority 15

			echo $prio > /sys/class/net/$DEVICE/default_priority
		fi
	;;
esac

